# chef-repo/cookbook_artifacts/MemeGen-9598b0e6d048774533289027e01a117be5858cbc/templates/myservice.service.erb
[Unit]
Description=Simplified Flask Service
After=syslog.target

[Service]
Type=simple
User=<%= @username %>
WorkingDirectory=/FeedbackApp-<%= @instance_name %>
ExecStart=/bin/bash -c 'sudo python3 -u /MemeGen/server.py <%= @port%>'
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.targeti